# inittab
# Begin /etc/inittab

id:3:initdefault:

rc::bootwait:/etc/rc.d/rc

su:1:wait:/usr/bin/bash

# End /etc/inittab